# roles/server/templates/Restart-Minesible.service
[Unit]
Description=Restart Minesible

[Service]
Type=oneshot
ExecStart=/opt/minecraft/backup.sh
WorkingDirectory=/opt/minecraft